# bin/rules.d/mega.rules
# Leonardo Nodes
SUBSYSTEM=="tty", SUBSYSTEMS=="usb", ENV{ID_MODEL}=="Arduino_Mega",  ATTRS{idProduct}=="0042", SYMLINK+="ttyON_MEGA", RUN+="/bin/stty -F /dev/ttyON_MEGA 57600"


SUBSYSTEMS=="usb", ATTRS{idProduct}=="0042", ATTRS{idVendor}=="2a03", MODE="0664", GROUP="dialout"
